# cluster/addons/haproxy-sni/images/supervisord.conf
[supervisorctl]
serverurl="unix:///var/run/supervisor.sock"

[supervisord]
nodaemon=true

[program:blbmonitor]
command=/bin/haproxy-monitor
startretries=100
autorestart=true
